# native-host/com.claude.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.app</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/umutbasilgan/BrowserExtension/native-host/venv/bin/python3</string>
        <string>/Users/umutbasilgan/BrowserExtension/native-host/claude-proxy.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/claude-proxy.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-proxy.stderr.log</string>
</dict>
</plist> 